# .crontab
* * * * * php /var/www/zonglist/cron.php /crons/index #_zong_cron
